# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm7

[program:nginx]
command=nginx
