# rootfs/services/php8/run
#!/usr/bin/env sh
/usr/sbin/php-fpm8
s6-svscanctl -qz /etc/s6.d
